# workspace/cron/cron.d/ufbackup-mediawiki
20 21 * * * root /usr/local/bin/ufbackup-mediawiki.sh
