# workspace/cron/cron.d/ufbackup-git
41 20 * * * root /usr/local/bin/ufbackup-git-repos.sh

