# workspace/cron/cron.d/ufcleanup-devopswiki
0 22 * * * root /usr/local/bin/ufcleanup-devopswiki.sh
