# workspace/cron/cron.d/ufcleanup-ufdist
10 22 * * * root /usr/local/bin/ufcleanup-ufdist.sh
