# workspace/cron/cron.d/ufcleanup-ee
30 01 * * * root /usr/local/bin/ufcleanup-ee.sh

