# workspace/cron/cron.d/ufcleanup-mail
35 01 * * * root /usr/local/bin/ufcleanup-mail.sh

