# workspace/cron/cron.d/ufcleanup-ufo-log
25 * * * * root /usr/local/bin/ufcleanup-ufo-log.sh
