# workspace/cron/cron.d/ufcleanup-ufo
15 01 * * * root /usr/local/bin/ufcleanup-ufo.sh

