# workspace/cron/cron.d/ufmaster-switch
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
55 5 * * 1 root /usr/local/bin/ufmaster-switch.sh auto >> /var/log/ufmaster-switch.log 2>&1
