# workspace/cron/cron.d/ufrestart-munin
22 01 * * * root /usr/sbin/service munin-node restart

