# workspace/cron/cron.d/xfs
0 5 * * * root /usr/sbin/xfs_fsr -t 1800
